# service/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:django]
directory=/app
command=/usr/local/bin/gunicorn --conf gunicorn.conf.py
stdout_logfile=/var/log/supervisor/django.log
stdout_logfile_maxbytes=5MB
redirect_stderr=true
